# scripts/switchbot-controller.service
[Unit]
Description=SwitchBot Controller WebApp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=t-tsukahara
WorkingDirectory=/home/t-tsukahara/switchbot-controller
EnvironmentFile=/home/t-tsukahara/switchbot-controller/.env
ExecStart=/home/t-tsukahara/switchbot-controller/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
